# nft-blackhole-reload.timer
[Unit]
Description=Update nft-blackhole lists

[Timer]
OnBootSec=6h
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
