# nft-blackhole-reload.timer
[Unit]
Description="Update blackhole lists"

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h

[Install]
WantedBy=multi-user.target
